# deploy/aktuba-bot.service
[Unit]
Description=Aktuba Telegram Bot
After=network.target

[Service]
Type=simple
User=REPLACE_WITH_SERVER_USER
WorkingDirectory=/opt/aktuba_bot
EnvironmentFile=/opt/aktuba_bot/.env
ExecStart=/opt/aktuba_bot/.venv/bin/python /opt/aktuba_bot/bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

